# service-units/rclone-2ndbrain-bisync.service
[Unit]
Description=RClone Bisync for 2ndBrain (workstation)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync gdrive-vault: %h/Documents/2ndBrain \
  --password-command "@@RCLONE_PASSWORD_CMD@@" \
  --resilient \
  --recover \
  --conflict-resolve newer
