# service-units/rclone-2ndbrain.service
[Unit]
Description=RClone Mount for 2ndBrain
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/rclone mount gdrive-vault: %h/Documents/2ndBrain \
  --password-command "@@RCLONE_PASSWORD_CMD@@" \
  --vfs-cache-mode full \
  --vfs-cache-max-age 24h \
  --vfs-cache-max-size 10G \
  --poll-interval 15s \
  --dir-cache-time 5s
ExecStop=/usr/bin/fusermount -u %h/Documents/2ndBrain
Restart=on-failure

[Install]
WantedBy=default.target
